# jobs/cf-redis-broker/templates/backup.cron.erb
5 * * * * /var/vcap/jobs/cf-redis-broker/bin/redis_hourly-backup 2>&1
